# build/sonar.conf
# Info
# /etc/init place it! Service sonar start
description "Sonar is the site visitors monitor"
author      "Stanislav WEB"

# Events
start on startup
stop on shutdown

# Automatically respawn
respawn
respawn limit 20 5

# Run the script!
# Note, in this example, if your PHP script (the socket) returns
# the string "ERROR", the daemon will stop itself.
script
    [ $(exec /usr/bin/php -f /path/to/cli command) = 'ERROR' ] && ( stop; exit 1; )
end script